# containers/crontabs/debug-minutely.cron
## Debug - once per minute
* * * * * /bin/bash /git-sync/git_mirror_sync.sh 2>&1 | tee -a /var/log/cron.log > /proc/1/fd/1
