# containers/crontabs/hourly.cron
@hourly /bin/bash /git-sync/git_mirror_sync.sh >&1 | tee -a /var/log/cron.hourly.log > /proc/1/fd/1
